# kafka/supervisor/kafka-create-topic.conf
[program:kafkaTopic]
command=/usr/bin/create-kafka-topic.sh
autostart=false
autorestart=false